# cherokee.upstart

start on pld.network-started and started random
stop on pld.shutdown-started

env SERVICE_web_server=y
export SERVICE_web_server

respawn

console output

script
	exec /usr/sbin/cherokee
end script

# vi: ft=upstart
